# glean/init/glean.openrc
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# $Id$

description="basic instance setup"

start() {
  %%GLEAN_SCRIPTS_DIR%%/glean-legacy.sh
  eend 0
}

depend() {
  need root localmount
  before net
}
